# profile.d/10-env-vars.sh
#! /usr/bin/env bash
# Personal environment variables
export EDITOR="code --wait"
export GITHUB_USER=""
export EMAIL=""
export NAME=""

# NPM
export NPM_CONFIG_INIT_AUTHOR_EMAIL="$EMAIL"
export NPM_CONFIG_INIT_AUTHOR_NAME="$NAME"
export NPM_CONFIG_INIT_LICENSE="Apache-2.0"
export NPM_CONFIG_INIT_VERSION="0.0.0"
